# wiki/cablecat-cleanup.service
[Unit]
Description=Clean up old CableCat Wiki cache files

[Service]
Type=oneshot
ExecStart=/usr/bin/cablecat-cleanup
